# wireguard-fs/etc/init.d/wg-quick
#!/sbin/openrc-run

description="Amnezia WG Quick Up"

depend() {
    need localmount
    need net
}

start_pre()
{
    /data/pre_up.sh
}

start() {
    for file in $(grep '\[Interface\]' /etc/amnezia/amneziawg/*.conf -l 2>/dev/null); do
        interface=$(basename $file .conf)
        ebegin "Starting Amnezia WG interface $interface"
        /usr/bin/awg-quick up $file
    done
    eend 0
}

post_start() {
    # Проверка, что интерфейс awg0 активен
    if ip addr show | grep -q "awg0"; then
        einfo "Applying iptables rule for awg0 interface"
        iptables -t nat -A POSTROUTING -o awg0 -j MASQUERADE
    else
        ewarn "awg0 interface is not active"
    fi
}

stop() {
    for file in $(grep '\[Interface\]' /etc/amnezia/amneziawg/*.conf -l 2>/dev/null); do
        interface=$(basename $file .conf)
        ebegin "Starting Amnezia WG interface $interface"
        /usr/bin/awg-quick down $file
    done
    eend 0
}
